# ProgressBot.service
[Unit]
Description=School Day Progress Bot
After=network-online.target

[Service]
Type=simple
Environment="LENGTH=20"
Environment="START_DATE=yyyy-mm-dd"
Environment="END_DATE=yyyy-mm-dd"
Environment="API_TOKEN=123456789:ABCdEFGHIJKLMnop1QrSt2Uv3Wx_yzaBCdef"
Environment="CHAT_ID=-1234567890"
Environment="WORKDIR=/path/to/ProgressBot"
ExecStart=/bin/bash /path/to/ProgressBot/main.sh
Restart=always
User=nobody

[Install]
WantedBy=multi-user.target
